# files/lede_generic/gl_b1300/etc/init.d/check_wireless
#!/bin/sh /etc/rc.common
# (C) 2018 minim.co

# When to start/stop the agent during the system init/deinit
START=05
boot() {
	[ -f /etc/config/wireless ] && 	touch /tmp/wireless_exists.txt
}
