# files/lede_generic/gl_ar300m/etc/init.d/set_defaults
#!/bin/sh /etc/rc.common
# (C) 2018 minim.co

# When to start/stop the agent during the system init/deinit
START=12
boot() {
  /sbin/set_defaults.sh
}
