# files/lede_generic/restart_config.init
#!/bin/sh /etc/rc.common
# (C) 2018 minim.co

# When to start 
START=11
boot() {
  /usr/bin/restart_config.sh 0
}
