# files/lede_generic/minim_os/etc/uci-defaults/31_dummy_wan
uci -q batch <<-EOF
	set network.dummywan='interface'
	set network.dummywan.device='eth1'
	set network.dummywan.proto='none'
