# files/lede_generic/minim_q14_openwifi/etc/uci-defaults/31_default_wireless
# If already configured via the cloud, exit
uci -q get minim.@unum[0].opmode && exit 0

# SSID
ssid=q14-`fw_printenv -n ethaddr | tr 'A-F' 'a-f' | awk 'BEGIN{FS=":"} {printf("%s%s\n", $5, $6)}'`
if [ -z "$ssid" ]; then
	exit 0
fi

# Key
key=`fw_printenv -n wifi_psk`
if [ -z "$key" ]; then
	exit 0
fi

#Set the default SSID and Key for both the radios
for devidx in 0 1 2
do
	uci -q batch <<-EOF
		set wireless.@wifi-iface[${devidx}].ssid=$ssid
		set wireless.@wifi-iface[${devidx}].encryption="sae-mixed"
		set wireless.@wifi-iface[${devidx}].key=$key
		set wireless.@wifi-iface[${devidx}].disabled=0
		set wireless.@wifi-iface[${devidx}].wps_pushbutton=1
		set wireless.@wifi-iface[${devidx}].wps_device_name=$ssid
		set wireless.@wifi-iface[${devidx}].wps_manufacturer="Motorola"
		set wireless.radio${devidx}.disabled=0
		set wireless.radio${devidx}.country=US
	EOF
done
uci -q commit wireless
exit 0
