# files/lede_generic/minim_q14_openwifi/etc/uci-defaults/33_default_password
# password
password=`fw_printenv -n wifi_psk`
if [ -z "$password" ]; then
	exit 0
fi

# Set the password
echo -e "$password\n$password" | passwd root

exit 0
