# files/lede_generic/minim_os/etc/hotplug.d/iface/99-platform
if [ "$INTERFACE" == "wan" ] && [ "$ACTION" == "ifup" ]; then
	/sbin/platform_boot.sh
fi
