# files/lede_generic/archer_c7_v4ap/etc/init.d/set_defaults
#!/bin/sh /etc/rc.common
# (C) 2019 minim.co

START=12
boot() {
  /sbin/set_defaults.sh
}
