# files/lede_generic/minim_os/etc/init.d/minim_done
#!/bin/sh /etc/rc.common
# (C) 2023 minim.co

# Use this place for platform specific scripts
START=99
boot() {
  /sbin/platform_boot.sh
}
