# files/lede_generic/tl_wdr3600_v1/etc/hotplug.d/iface/30-set-lan-mac
if [ "$INTERFACE" == "lan" ] &&
	[ "$ACTION" == "ifup" ]; then
	/sbin/set_lan_iface.sh
fi
